# guzald.service
[Unit]
Descrption=Guzal-Aida daemon

[Service]
Type=simple
ExecStart=/usr/bin/daemon

[Install]
WantedBy=multi-user.target